# crontab
# run script every 1st in the month

0 0 1 * * /root/nextjs-socket.io/scripts/renew.sh